# runner.vm/files/usr/lib/systemd/system/limit.timer
[Unit]
Description=Limit uptime to this

[Timer]
OnBootSec=60m
Unit=poweroff.target

[Install]
WantedBy=timers.target
